# testdata/net/dhcp/docker/supervisord.conf
[supervisord]
nodaemon=true

[program:dhcpd4]
command=/usr/sbin/dhcpd -4 -f -d -cf /root/dhcpd_4.conf -lf /root/dhcpd.leases 
autostart=false

[program:dhcpd6]
command=/usr/sbin/dhcpd -6 -f -d -cf /root/dhcpd_6.conf -lf /root/dhcpd.leases 
autostart=false

[unix_http_server]
file=/tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock
